# crontab
10 */1 * * * /bin/scoober-sync >> /var/log/scoober-sync.log
